# 80-rfid.rules
SUBSYSTEMS=="usb" ATTRS{idVendor}=="05fe" ATTRS{idProduct}=="1010" RUN+="/etc/init.d/xbmc-rfid-music start" 
